# launchd/com.fracto.evidence-pack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fracto.evidence-pack</string>
  <key>ProgramArguments</key>
  <array>
    <string>$HOME/myprojects/mcc/local-isms-starter-kit/.venv/bin/python</string>
    <string>$HOME/myprojects/mcc/local-isms-starter-kit/evidence_pack/evidence_pack_local.py</string>
    <string>--run-once</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key><integer>1</integer>
    <key>Hour</key><integer>5</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>$HOME/Library/Logs/fracto.evidence-pack.out</string>
  <key>StandardErrorPath</key><string>$HOME/Library/Logs/fracto.evidence-pack.err</string>
  <key>RunAtLoad</key><true/>
</dict>
</plist>
